# packages/deepshit-network/debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        echo "Installing Reticulum and Sideband..."
        if ! pip3 install --quiet reticulum sideband; then
            echo "Failed to install Reticulum dependencies." >&2
            exit 1
        fi

        install -d /etc/reticulum
        install -d /var/lib/deepshit/reticulum

        if [ ! -f /var/lib/deepshit/reticulum/identity.id ]; then
            if ! /usr/bin/deep-generate-identity; then
                echo "Failed to generate DeepShitNet identity." >&2
                exit 1
            fi
        fi

        if command -v systemctl >/dev/null 2>&1; then
            systemctl daemon-reload || true
            systemctl enable deepshit-reticulum.service >/dev/null 2>&1 || true
            systemctl restart deepshit-reticulum.service >/dev/null 2>&1 || true
        fi

        echo "DeepShitNet networking setup complete."
        ;;
esac

exit 0
