# packages/deepshit-network/src/etc/systemd/system/deepshit-reticulum.service
[Unit]
Description=DeepShitNet Reticulum Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/python3 -m RNS.Utilities.Reticulum
Restart=on-failure

[Install]
WantedBy=multi-user.target
